# chef/cookbooks/xvfb/files/default/xvfb.conf
# Upstart file for Xvfb

start on runlevel [2345]
stop on runlevel [06]

respawn

exec Xvfb :0 -screen 0 1024x768x24