# chef/cookbooks/jmeter/files/default/jmeter-server.conf
# jmeter-server startup

description "jmeter-server startup"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec jmeter-server