# chef/cookbooks/reportgrid/templates/default/blueeyes_service.init.conf.erb
# <%= @service_name %>-<%= @version %> - ReportGrid <%= @service_name %> service (<%= @version %>)
#
#

description "<%= @service_name %>-<%= @version %> - ReportGrid <%= @service_name %> service (<%= @version %>)"

start on runlevel [2345]
stop on runlevel [06]

respawn

pre-start script
  [ -d /var/run/reportgrid ] || install -m 755 -o reportgrid -g reportgrid -d /var/run/reportgrid
end script

script
  ulimit -n 65535
  if [ -f /etc/default/<%= @service_name %>-<%= @version %> ]; then . /etc/default/<%= @service_name %>-<%= @version %>; fi
  su reportgrid -s "/bin/bash" -c "/usr/bin/java $JAVA_OPTS -jar /usr/share/java/<%= @service_name %>-<%= @version %>.jar --configFile /etc/reportgrid/<%= @service_name %>-<%= @version %>.conf"
end script

post-start script
  echo `status <%= @service_name %>-<%= @version %> | egrep -o '[0-9]+$' | head -n1` > /var/run/reportgrid/<%= @service_name %>-<%= @version %>.pid
end script
